# config/elasticsearch/n0/elasticsearchnode0.conf

[program:elasticsearchnode0]
command=/opt/elasticsearch6/bin/elasticsearch
    -Epath.conf=/usr/local/etc/elasticsearch/n0
autostart=true
autorestart=true
startsecs=5
startretries=3
user=_elastic
environment=
    ES_JVM_OPTIONS="/usr/local/etc/elasticsearch/n0/jvm.options",
    JAVA_HOME="/opt/jdk8"
stderr_logfile=/var/log/supervisor/%(program_name)s.err.log
stdout_logfile=/var/log/supervisor/%(program_name)s.out.log


